# scripts/supervisor/aura-bootstrap.conf
[program:aura-bootstrap]
command=/usr/bin/env bash /path/to/your/repo/scripts/bootstrap-aura-idtoken.sh
directory=/path/to/your/repo
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/var/log/aura-bootstrap.out.log
stderr_logfile=/var/log/aura-bootstrap.err.log
user=youruser
environment=COMPOSE_FILE="/path/to/your/repo/docker-compose.yml"
